# kernse/deploy/vps/systemd/kernse-admin-v2.service
# Kernse — Console plateforme v2 (admin.kernse.fr)
# Accès strictement réservé aux superadmins plateforme (auth session + 2FA TOTP).
# Recommandation nginx : IP allowlist en amont.

[Unit]
Description=Kernse Admin v2 (admin.kernse.fr)
After=network.target
StartLimitIntervalSec=200
StartLimitBurst=5

[Service]
Type=simple
User=kernse
Group=kernse
WorkingDirectory=/home/kernse/platform

EnvironmentFile=/home/kernse/env/kernse-admin-v2.env
Environment="KERNSE_ENV=v2"
Environment="KERNSE_ADMIN_PORT=8102"
Environment="KERNSE_PLATFORM_DB_PATH=/home/kernse/platform/kernse/admin/data/platform.db"

ExecStart=/home/kernse/venv/bin/uvicorn kernse.admin.main:app --host 127.0.0.1 --port 8102 --workers 1

Restart=on-failure
RestartSec=5

# NoNewPrivileges DOIT être false : kernse-admin appelle sudo (via la
# whitelist /etc/sudoers.d/kernse-admin) pour lancer provision_client.sh
# et promote_client.sh. Le flag à true empêcherait sudo d'élever les
# privilèges et casserait le provisioning et la promotion.
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=false
# Admin peut écrire dans /home/kernse (DB, backups, scripts) ET dans les
# chemins /etc/systemd, /etc/nginx et /etc/letsencrypt via sudo (script
# provision_client.sh). ProtectSystem=strict rend le reste de / en RO même
# via sudo — donc chaque chemin utilisé par provision_client.sh doit être
# explicitement whitelisté ici.
ReadWritePaths=/home/kernse /etc/systemd/system /etc/nginx/sites-available /etc/nginx/sites-enabled /var/www/certbot /etc/letsencrypt
PrivateTmp=true

[Install]
WantedBy=multi-user.target
